# docker/debian8/ed/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:postgresql]
command=/usr/sbin/service postgresql start
